# deploy/nekko-collector.service
# NEKKO F2 — unit de systemd para el colector (ADR-0007)
# Instalación: copiar a /etc/systemd/system/, luego
#   systemctl daemon-reload && systemctl enable --now nekko-collector
[Unit]
Description=NEKKO F2 live collector (Polymarket/Binance, read-only)
After=network-online.target
Wants=network-online.target

[Service]
User=nekko
WorkingDirectory=/home/nekko/9-NEKKO
ExecStart=/home/nekko/9-NEKKO/.venv/bin/python -m collector.live_collector
Restart=always
RestartSec=15
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
